# packages/unv.deploy/unv.deploy-0.2.10-py3-none-any.whl/unv/deploy/components/app/app.service
[Unit]
Description={{ description }} [{{ instance }}]
After=network.target

[Service]
Type=simple
User={{ settings.user }}
Group={{ settings.user }}
ExecStart={{ settings.bin.format(user=user, public_ip=public_ip, private_ip=private_ip, settings=settings, instance=instance) }}
WorkingDirectory={{ settings.home_abs }}
LimitNOFILE={{ limit_nofile }}

Environment=PYTHONPATH={{ settings.home_abs }}
Environment=SETTINGS={{ settings.module }}

[Install]
WantedBy=multi-user.target
